# resources/s6-rc.d/celery/run
#!/command/execlineb -P

with-contenv
celery -A config.celery worker -l INFO
