# resources/s6-rc.d/redis/run
#!/command/execlineb -P

cd /var/lib/redis
with-contenv
s6-setuidgid redis
redis-server
