# resources/s6-rc.d/setup1/up
# execline

cd /opt/services/djangoapp/requirements
foreground { touch requirements.txt }
pip install -r requirements.txt
