# resources/s6-rc.d/djangoapp/finish
#!/command/execlineb -S2

s6-svscanctl -t /run/service
