# resources/s6-rc.d/djangoapp/run
#!/command/execlineb -P

cd /opt/services/djangoapp/src
with-contenv
gunicorn -c config/gunicorn.conf.py config.wsgi:application
