# resources/s6-rc.d/nginx/run
#!/command/execlineb -P

cd /opt/services/djangoapp
foreground { chown www-data:www-data media }
with-contenv
nginx -g "daemon off;"
