# resources/s6-rc.d/setup2/up
# execline

cd /opt/services/djangoapp/src
with-contenv
python3 manage.py collectstatic --noinput
