# resources/s6-rc.d/setup2/up
#!/command/execlineb -P

cd /opt/services/djangoapp/src
with-contenv
python3 manage.py collectstatic --noinput
